# home/Library/LaunchAgents/com.oho.cntlm.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.oho.cntlm.daemon</string>
    <key>ProgramArguments</key>
    <array>
      <string>/opt/local/sbin/cntlm</string>
      <string>-v</string>
      <string>-c</string>
      <string>/Users/josephwinston/Tools/x86_64-apple-darwin15.2.0/etc/cntlm.conf</string>
    </array>
    <key>KeepAlive</key>
    <false/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/josephwinston/Library/Logs/cntlm.log</string>
    <key>StandardOutPath</key>
    <string>/Users/josephwinston/Library/Logs/cntlm.log</string>
  </dict>
</plist>

